# etc/udev/rules.d/99-intel-rapl.rules
SUBSYSTEM=="powercap", ACTION=="add", RUN+="/usr/bin/systemd-tmpfiles --create /etc/tmpfiles.d/intel-rapl.conf"
