# etc/udev/rules.d/turboboost.rules
RUN+="/bin/chgrp turboboost /sys/devices/system/cpu/intel_pstate/no_turbo"
RUN+="/bin/chmod g+w /sys/devices/system/cpu/intel_pstate/no_turbo"
